# templates/dehydrated-renew.timer.j2
[Unit]
Description=Timer for Dehydrated certificate renewal

[Timer]
OnCalendar={{ dehydrated_timer_renew }}
Persistent=true

[Install]
WantedBy=timers.target
