# templates/dehydrated-check.timer.j2
[Unit]
Description=Timer for Dehydrated certificate expiry check

[Timer]
OnCalendar={{ dehydrated_timer_check }}
Persistent=true

[Install]
WantedBy=timers.target
